# rootfs/etc/xdg/autostart/satdump-ui.desktop
[Desktop Entry]
Type=Application
Name=SatSumo
Exec=/usr/bin/satdump-ui
Terminal=false